# dream-server/scripts/systemd/dream-mdns.service
[Unit]
Description=Dream Server mDNS Announcer
After=network-online.target avahi-daemon.service
Wants=network-online.target

[Service]
Type=simple
User=__INSTALL_USER__
ExecStart=__PYTHON3__ __INSTALL_DIR__/bin/dream-mdns.py
WorkingDirectory=__INSTALL_DIR__
Environment=DREAM_INSTALL_DIR=__INSTALL_DIR__
Environment=HOME=__HOME__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
